# files/system/shared/usr/lib/systemd/system/amd-gpu-high-power.service
[Unit]
Description=Set amd gpu governor to performance
[Service]
ExecStart=/usr/libexec/amd-gpu-high-power
[Install]
WantedBy=multi-user.target